# system_files/shared/usr/lib/systemd/system/bamos-cache-clean.timer
[Unit]
Description=Weekly BamOS Cache Cleanup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
