# system_files/shared/usr/lib/systemd/system/flatpak-repair.timer
[Unit]
Description=Weekly Flatpak repair and metadata refresh
Documentation=https://github.com/quocnho/bamos

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
