# system_files/shared/usr/lib/systemd/system/flatpak-cleanup.timer
[Unit]
Description=Weekly Flatpak Cleanup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
