# system_files/shared/usr/lib/systemd/system/podman-prune.service
[Unit]
Description=Podman Container and Image Prune

[Service]
Type=oneshot
ExecStart=/usr/bin/podman system prune -f
IOSchedulingClass=idle
Nice=19
